# uflinux/files/99-bindfs.rules
ACTION=="add", ENV{ID_FS_USAGE}=="filesystem", RUN+="/system/fswrap.sh %E{DEVNAME}"
